# infinitetrade.service
[Unit]
Description=InfiniteTrade Bot
After=network.target

[Service]
User=gwyoon
WorkingDirectory=/home/gwyoon/auto_infinitetrade
ExecStart=/home/gwyoon/auto_infinitetrade/.venv/bin/python -m src.main
Restart=always
RestartSec=10
Environment=TZ=US/Eastern

[Install]
WantedBy=multi-user.target
